# borg-backup.service
[Unit]
Description=BorgBackup

[Service]
Type=simple
ExecStart=/root/borg_backup.sh