# var-backup.mount
[Unit]
Description=Mount backup volume sdb

[Mount]
What=/dev/sdb
Where=/var/backup
Type=xfs
Options=defaults

[Install]
WantedBy=multi-user.target